# specs/pkgs/man-db/etc/cron.daily/man-db
#!/bin/ash

/bin/nice -n 19 /bin/ionice -c 2 -n 7 /usr/bin/mandb -q
