# specs/pkgs/pciutils/etc/cron.daily/pci.ids
#!/bin/ash

exec /usr/sbin/update-pciids -q
