# specs/pkgs/findutils/etc/cron.daily/updatedb
#!/bin/ash

exec /bin/nice -n 19 /bin/ionice -c 2 -n 7 /usr/bin/updatedb
