# specs/pkgs/busybox/etc/inittab
::sysinit:-/etc/sv/1
::shutdown:-/etc/sv/3
::ctrlaltdel:/sbin/reboot
::restart:/sbin/init
::respawn:/etc/sv/2
