# specs/images/iso_initrd/etc/profile.d/cdrom-bindirs
export PATH="$PATH:/mnt/cdrom/bin:/mnt/cdrom/sbin"
